# cron/tasks/lmms/systemd/php-watcher.path
[Path]
PathModified=/home/deploy/lmms.io/public/locale

[Install]
WantedBy=multi-user.target

