# cron/tasks/lmms/systemd/php-watcher.service
[Unit]
Description=PHP FPM Restart Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart php7.0-fpm.service

